# fantasy-baseball.service
[Unit]
Description=Fantasy Baseball Dashboard
After=network.target
Wants=network.target

[Service]
Type=simple
User=YOUR_USER
Group=YOUR_USER
WorkingDirectory=/path/to/FantasyBaseball
Environment=PATH=/path/to/FantasyBaseball/venv/bin
Environment=PYTHONPATH=/path/to/FantasyBaseball/src
ExecStart=/bin/bash /path/to/FantasyBaseball/scripts/start_server.sh
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=fantasy-baseball

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/home/chandlergray/FantasyBaseball

[Install]
WantedBy=multi-user.target
